# deployment/omni-alpha.service
[Unit]
Description=Omni Alpha 5.0 - World-Class Trading System
After=network.target
Wants=network.target

[Service]
Type=exec
User=omni
Group=omni
WorkingDirectory=/opt/omni-alpha
Environment=PYTHONPATH=/opt/omni-alpha
Environment=ENV=production
Environment=LOG_LEVEL=INFO
ExecStart=/opt/omni-alpha/venv/bin/gunicorn src.app:app \
    --worker-class uvicorn.workers.UvicornWorker \
    --workers 4 \
    --bind 0.0.0.0:8000 \
    --max-requests 10000 \
    --max-requests-jitter 1000 \
    --timeout 60 \
    --keepalive 5 \
    --log-level info \
    --access-logfile - \
    --error-logfile -
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=omni-alpha

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/omni-alpha/logs /opt/omni-alpha/data
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
